# news-update.service
[Unit]
Description=News Trend Analysis - Daily Update
Wants=news-update.timer

[Service]
Type=oneshot
WorkingDirectory=/home/David/Documents/GithubRepos/news-trend-analysis
ExecStart=/home/David/Documents/GithubRepos/news-trend-analysis/.venv/bin/python daily_update.py --max-results 50
StandardOutput=append:/home/David/Documents/GithubRepos/news-trend-analysis/logs/systemd_update.log
StandardError=append:/home/David/Documents/GithubRepos/news-trend-analysis/logs/systemd_update.log

[Install]
WantedBy=multi-user.target
